# Linux/mc_atm10.service
[Unit]
Description=Minecraft Server ATM10
After=network.target

[Service]
Type=forking
User=mc
Group=mc
WorkingDirectory=/home/mc/ATM10

ExecStart=/usr/bin/screen -dmS MC_ATM10 /home/mc/ATM10/startserver.sh
ExecStop=/usr/bin/screen -S MC_ATM10 -X stuff "stop$(printf \\r)"
ExecStopPost=/bin/sleep 10

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
