# Linux/mc_atm10.service
[Unit]
Description=Minecraft Server ATM10
After=network.target

[Service]
User=mc
Group=mc
WorkingDirectory=/home/mc/ATM10
Environment="SERVERNAME=ATM10"
ExecStart=/bin/bash ./startserver.sh nogui
ExecStop=/usr/bin/su - mc -c "screen -S MC_ATM10 -p 0 -X stuff 'say Server stopping...'; screen -S MC_ATM10 -p 0 -X stuff 'stop^M'"
Restart=on-failure
RestartSec=10
KillMode=control-group
TimeoutStopSec=120
StartLimitInterval=600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
